# supervisord.conf
; =============================================================================
; supervisord.conf — Process manager for Hermes Suit all-in-one container
; Manages three services:
;   1. hermes-gateway   (port 8642) — Agent gateway (CLI, Telegram, etc.)
;   2. hermes-dashboard (port 9119) — Monitoring dashboard (built-in)
;   3. hermes-webui     (port 8787) — Browser chat interface
; =============================================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=10MB
logfile_backups=3
childlogdir=/tmp

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ---------------------------------------------------------------------------
; Service 1: Hermes Agent Gateway
; The main agent process — handles CLI, Telegram, cron, webhook, etc.
; ---------------------------------------------------------------------------
[program:hermes-gateway]
command=/opt/hermes/.venv/bin/hermes gateway run
directory=/opt/hermes
autostart=true
autorestart=true
priority=10
environment=HERMES_HOME="/opt/data",HERMES_DATA_DIR="/opt/data",HOME="/opt/data",PATH="/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

; ---------------------------------------------------------------------------
; Service 2: Hermes Dashboard
; Built-in analytics/monitoring dashboard (part of hermes-agent)
; ---------------------------------------------------------------------------
[program:hermes-dashboard]
command=/opt/hermes/.venv/bin/hermes dashboard --host 0.0.0.0 --port 9119 --insecure --no-open
directory=/opt/hermes
autostart=true
autorestart=true
priority=20
startsecs=5
environment=HERMES_HOME="/opt/data",HERMES_DATA_DIR="/opt/data",HOME="/opt/data",PATH="/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

; ---------------------------------------------------------------------------
; Service 3: Hermes WebUI
; Browser-based chat interface for interacting with the agent
; ---------------------------------------------------------------------------
[program:hermes-webui]
command=/opt/hermes-webui/venv/bin/python /opt/hermes-webui/server.py
directory=/opt/hermes-webui
autostart=true
autorestart=true
priority=30
startsecs=5
environment=HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_PORT="8787",HERMES_WEBUI_STATE_DIR="/opt/data/webui",HERMES_WEBUI_DEFAULT_WORKSPACE="/workspace",HERMES_HOME="/opt/data",HOME="/opt/data",PATH="/opt/hermes-webui/venv/bin:/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
